# scripts/macos/com.isl-assist-ai.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.isl-assist-ai.backend</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd "__REPO_ROOT__/backend" &amp;&amp; "__PYTHON__" -m uvicorn app.main:app --host 127.0.0.1 --port 8000</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/isl-assist-ai-backend.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/isl-assist-ai-backend.log</string>
  </dict>
</plist>
